# supervisor/conf.d/supervisord-logs.conf
[program:php7-fpm-log]
command=tail -f /var/log/php7-fpm.log
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
