# supervisor/conf.d/supervisord-php5-fpm.conf
[program:php5-fpm]
command=/etc/supervisor/commands/start-php5-fpm.sh
priority=999
username=www-data
autostart=true
autorestart=true
redirect_stderr=true
