# supervisor/conf.d/supervisord-logs.conf
[program:php5-fpm-log]
command=tail -f /var/log/php5-fpm.log
redirect_stderr=true
