# supervisor/conf.d/supervisord-php7-fpm.conf
[program:php7-fpm]
command=/etc/supervisor/commands/start-php7-fpm.sh
priority=999
username=www-data
autostart=true
autorestart=true
redirect_stderr=true
