# supervisor/conf.d/supervisord-logs.conf
[program:php7-fpm-log]
command=tail -f /var/log/php7-fpm.log
redirect_stderr=true
